# image/DEBIAN/postinst
#!/bin/bash

systemctl daemon-reload
systemctl -f enable cjdns
systemctl restart cjdns
